# ops/systemd/orchestra-autopilot.service
[Unit]
Description=Stellcodex Orchestra Autopilot Stack
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/root/workspace/ops/orchestra
ExecStart=/root/workspace/ops/orchestra/autopilot.sh start /root/workspace/_jobs
ExecStop=/root/workspace/ops/orchestra/autopilot.sh stop /root/workspace/_jobs
TimeoutStartSec=180
TimeoutStopSec=60
Environment=TZ=Europe/Istanbul

[Install]
WantedBy=multi-user.target
