# ops/systemd/stellcodex-tmp-cleanup.service
[Unit]
Description=STELLCODEX artifact cleanup (>24h)

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/root/workspace/scripts/stellcodex_artifact_cleanup.sh
