# ops/systemd/stellcodex-autopilot-weekly.service
[Unit]
Description=STELLCODEX weekly integrity autopilot check
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/workspace
ExecStart=/root/workspace/scripts/stellcodex_autopilot.sh weekly --archive
