# ops/systemd/stellcodex-247-watchdog.service
[Unit]
Description=Stellcodex 7/24 Watchdog Check
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/root/workspace/ops/orchestra
ExecStart=/root/workspace/ops/orchestra/stellcodex_watchdog_check.sh
